# opensovd-core/deploy/pi/systemd/opensovd-cert-rotation.timer
[Unit]
Description=Run OpenSOVD certificate rotation sweep daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=opensovd-cert-rotation.service

[Install]
WantedBy=timers.target
