# opensovd-core/deploy/pi/systemd/sovd-main.service
# SPDX-License-Identifier: Apache-2.0
# SPDX-FileCopyrightText: 2026 The Contributors to Eclipse OpenSOVD (see CONTRIBUTORS)
#
# See the NOTICE file(s) distributed with this work for additional
# information regarding copyright ownership.
#
# Phase 5 Line A D1 — sovd-main systemd unit on the Raspberry Pi bench.
#
# Binds 127.0.0.1:21002 and relies on the observer nginx front end for
# any bench-LAN access. This matches opensovd-pi.toml, the Phase 5
# deploy verification, and the current sovd-main HTTP safety guard
# which rejects plain HTTP on non-loopback listeners unless the binary
# is built with the insecure-http-fallback feature.

[Unit]
Description=Taktflow OpenSOVD sovd-main (Phase 5 Line A D1 Pi deploy)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=taktflow-pi
Group=taktflow-pi
WorkingDirectory=/opt/taktflow/sovd-main
ExecStart=/opt/taktflow/sovd-main/sovd-main \
    --config-file /opt/taktflow/sovd-main/opensovd.toml \
    --listen-address 127.0.0.1 \
    --listen-port 21002
Restart=on-failure
RestartSec=5
# Keep logs compact on the Pi (SD card lifetime).
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
