# opensovd-core/deploy/pi/systemd/taktflow-can-doip-proxy.service
# SPDX-License-Identifier: Apache-2.0
# SPDX-FileCopyrightText: 2026 The Contributors to Eclipse OpenSOVD (see CONTRIBUTORS)
#
# See the NOTICE file(s) distributed with this work for additional
# information regarding copyright ownership.
#
# Phase 5 Line A D1 — Taktflow CAN-to-DoIP proxy systemd unit.
#
# ============================================================
# CONFLICT NOTE: :13400 ecu-sim vs :13401 proxy
# ============================================================
# The upstream CDA ecu-sim simulator listens on :13400 (see
# ecu-sim.service). The Phase 2 Line B CAN-to-DoIP proxy is
# configured in this unit to bind :13401 (per phase-5-line-a.md port
# plan). That way the two can coexist on the Pi without a socket
# collision.
#
# IF the proxy is ever re-pointed at :13400 (e.g. to drop-in replace
# ecu-sim with a real-hardware proxy), this unit's `Conflicts=` line
# below guarantees systemd will stop ecu-sim.service first. Do NOT
# remove the Conflicts= line without also re-verifying that both
# services no longer share a listening port.
#
# The Phase 2 Line B proxy binary is produced by the Line B repo
# (taktflow-embedded-production, posix build). This unit assumes the
# binary is installed at /opt/taktflow/proxy/taktflow-can-doip-proxy.
# The deploy script (phase5-full-stack.sh) rsyncs it only if the
# source artifact resolves on the dev host; otherwise the proxy is
# skipped and this unit stays disabled (D1 green still holds on
# sovd-main alone — see README-phase5.md).

[Unit]
Description=Taktflow CAN-to-DoIP proxy (Phase 2 Line B, :13401)
After=network-online.target
Wants=network-online.target
# Hard-enforce the port-plan invariant: ecu-sim and proxy never both
# bind :13400. If a future config moves the proxy to :13400, systemd
# will stop ecu-sim first.
Conflicts=ecu-sim.service

[Service]
Type=simple
User=taktflow-pi
Group=taktflow-pi
WorkingDirectory=/opt/taktflow/proxy
ExecStart=/opt/taktflow/proxy/taktflow-can-doip-proxy \
    --config-file /opt/taktflow/proxy/proxy.toml \
    --listen-address 0.0.0.0 \
    --listen-port 13401 \
    --can-interface can0
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
