# opensovd-core/deploy/pi/ecu-sim.service
# SPDX-License-Identifier: Apache-2.0
# SPDX-FileCopyrightText: 2026 The Contributors to Eclipse OpenSOVD (see CONTRIBUTORS)
#
# See the NOTICE file(s) distributed with this work for additional
# information regarding copyright ownership.
#
# This program and the accompanying materials are made available under the
# terms of the Apache License Version 2.0 which is available at
# https://www.apache.org/licenses/LICENSE-2.0

[Unit]
Description=Taktflow CDA ECU Simulator (upstream Kotlin)
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=simple
ExecStartPre=-/usr/bin/docker rm -f ecu-sim
# Host networking is required so DoIP VAM broadcasts originate from the Pi's
# real network interface (not the 172.17.0.2 docker bridge address) and reach
# the dev machine on the same /24.
#
# The Pi has two interfaces: eth0 (203.0.113.0/24 — wired, not shared with the
# dev machine) and wlan0 (192.0.2.0/24 — shared with the dev machine). We
# force ecu-sim to bind wlan0 so VAM broadcasts reach the dev machine.
ExecStart=/usr/bin/docker run --rm --name ecu-sim --network host -e SIM_NETWORK_INTERFACE=wlan0 taktflow/cda-ecu-sim:latest
ExecStop=/usr/bin/docker stop ecu-sim
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
